# linux/monitor-batt.service
[Unit]
Description=Monitor-Batt — Daemon de monitoramento de bateria
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/monitor-batt/monitor_batt.py
Restart=on-failure
RestartSec=10
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target
